# common/supervisor/etc/supervisord.conf.d/health.conf
[program:health]
user=stellar
directory=/opt/stellar/health
command=/opt/stellar/health/bin/start
autostart=true
autorestart=true
priority=60
redirect_stderr=true