# common/supervisor/etc/supervisord.conf.d/nginx.conf
[program:nginx]
user=www-data
directory=/opt/stellar/nginx
command=/opt/stellar/nginx/bin/start
autostart=true
autorestart=true
priority=50
redirect_stderr=true
