# common/supervisor/etc/supervisord.conf.d/readiness.conf
[program:readiness]
user=stellar
directory=/opt/stellar/readiness
command=/opt/stellar/readiness/bin/start
autostart=true
autorestart=true
priority=60
redirect_stderr=true