# common/supervisor/etc/supervisord.conf.d/soroban-rpc.conf
[program:soroban-rpc]
user=stellar
directory=/opt/stellar/soroban-rpc
command=/opt/stellar/soroban-rpc/bin/start
autostart=false
startretries=50
autorestart=true
priority=60
redirect_stderr=true
