# common/supervisor/etc/supervisord.conf.d/stellar-core.conf
[program:stellar-core]
user=stellar
directory=/opt/stellar/core
command=/opt/stellar/core/bin/start
autostart=false
autorestart=true
priority=20
