# common/supervisor/etc/supervisord.conf.d/stellar-rpc.conf
[program:stellar-rpc]
user=stellar
directory=/opt/stellar/stellar-rpc
command=/opt/stellar/stellar-rpc/bin/start
autostart=false
startretries=50
autorestart=true
priority=60
redirect_stderr=true
