# common/supervisor/etc/supervisord.conf.d/stellar-lab.conf
[program:stellar-lab]
user=stellar
directory=/opt/stellar/lab
command=/opt/stellar/lab/bin/start
autostart=false
startretries=50
autorestart=true
priority=60
redirect_stderr=true
