# common/supervisor/etc/supervisord.conf.d/stellar-core-watcher.conf
[program:stellar-core-watcher]
user=stellar
directory=/opt/stellar/core-watcher
command=/opt/stellar/core-watcher/bin/start
autostart=false
autorestart=true
priority=21
