# local/supervisor/etc/supervisord.conf.d/galexie.conf
[program:galexie]
user=stellar
directory=/opt/stellar/galexie
command=/opt/stellar/galexie/bin/start
autostart=false
startretries=50
autorestart=true
priority=70
redirect_stderr=true
