# standalone/supervisor/etc/supervisord.conf.d/friendbot.conf
[program:friendbot]
user=stellar
directory=/opt/stellar/friendbot
command=/opt/stellar/friendbot/bin/start
autostart=true
autorestart=true
priority=40
redirect_stderr=true
